# systemd/ubuntu-auto-update-dashboard.service
[Unit]
Description=Ubuntu Auto-Update Dashboard
Documentation=https://github.com/patel5d2/ubuntu-auto-update
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# Adjust this to where you place the repo and venv
WorkingDirectory=/opt/ubuntu-auto-update
EnvironmentFile=/etc/ubuntu-auto-update-dashboard.env
# Requires a venv with dashboard requirements installed
ExecStart=/opt/ubuntu-auto-update/.venv/bin/python -m uvicorn dashboard.app:app --host 0.0.0.0 --port ${PORT}
Restart=on-failure
RestartSec=3

# Security hardening (adjust if needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
