# systemd/ubuntu-auto-update.service
[Unit]
Description=Ubuntu Auto-Update Service
Documentation=https://github.com/yourusername/ubuntu-auto-update
Wants=network-online.target
After=network-online.target
ConditionPathExists=/usr/local/bin/update.sh

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/bin/update.sh --quiet
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ubuntu-auto-update

# Security settings
NoNewPrivileges=yes
ProtectSystem=false
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=no
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectKernelModules=no
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes

# Resource limits
TimeoutStartSec=3600
TimeoutStopSec=30
MemoryMax=512M
TasksMax=100

# Environment
Environment="DEBIAN_FRONTEND=noninteractive"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target